# phantom-screenshot.path
[Unit]
Description=Phantom screenshot request watcher

[Path]
PathExists=/tmp/phantom_ss_request
Unit=phantom-screenshot.service

[Install]
WantedBy=default.target
